# deployment/supervisor/conf.d/application.conf
[program:application]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan serve --host=0.0.0.0 --port=9000
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/var/www/storage/logs/application.log
stdout_logfile_maxbytes=1024
stdout_logfile_backups=0
