# deployment/supervisor/conf.d/horizon.conf
[program:horizon]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan horizon
autostart=true
autorestart=true
numprocs=3
stdout_logfile=/var/www/storage/logs/horizon.log
stdout_logfile_maxbytes=1024
stdout_logfile_backups=0
